# deploy/trading-bot.service
[Unit]
Description=Trading Bot (SMC+PA)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/trading-bot
ExecStart=/opt/trading-bot/.venv/bin/python main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment
EnvironmentFile=/opt/trading-bot/.env

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
# .tmp holds the Shioaji simulation runtime files (SJ_LOG_PATH, etc.)
# and the local-paper-runner state JSON; the service crashes on
# startup without write access here.
ReadWritePaths=/opt/trading-bot/logs /opt/trading-bot/data /opt/trading-bot/reports /opt/trading-bot/.tmp
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-bot

[Install]
WantedBy=multi-user.target
